# etc/supervisord/supervisord.conf
[supervisord]
nodaemon = true
user = root
# loglevel = debug

[program:sshd]
command = /usr/sbin/sshd -D

[program:rstudio-server]
command = /etc/init.d/rstudio-server start


